# deploy/systemd/cert-watch.service
[Unit]
Description=cert-watch — TLS certificate expiry tracker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cert-watch
Group=cert-watch
Environment=CERT_WATCH_DATA_DIR=/var/lib/cert-watch
# Scanning options
# Environment=CERT_WATCH_ALLOW_PRIVATE_IPS=1
# Secrets: set fixed values in production to survive restarts
# Environment=CERT_WATCH_CSRF_SECRET=change-me-in-production
# Environment=CERT_WATCH_AUTH_SECRET=change-me-in-production
ExecStart=/opt/cert-watch/venv/bin/cert-watch
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/cert-watch
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
